# deploy/comroster.service
# ⚠️ MODÈLE DE RÉFÉRENCE — CE FICHIER N'EST PAS CELUI QUI TOURNE SUR LE BOÎTIER.
#
# Il décrit une installation MANUELLE « serveur classique » : utilisateur dédié
# `comroster`, application dans /opt/comroster, derrière Nginx (deploy/nginx.conf).
#
# Sur un Raspberry Pi, c'est deploy/setup-pi.sh qui GÉNÈRE le vrai unit dans
# /etc/systemd/system/comroster.service, avec des valeurs différentes :
#   • User      = l'utilisateur qui a lancé sudo (${SUDO_USER:-pi}), pas `comroster`
#   • Répertoire = le dépôt cloné, pas /opt/comroster
#   • Secrets   = EnvironmentFile=/etc/comroster.env, pas des Environment= en clair
# Pour connaître l'utilisateur réel d'un boîtier : `systemctl show -p User --value comroster`
#
# NB : volontairement PAS de NoNewPrivileges ici — il EMPÊCHE sudo de fonctionner, donc
# les boutons « Redémarrer le boîtier » et « Appliquer maintenant » (sudo strictement
# limité) échoueraient. Ne le réactive que si tu renonces à ces fonctions.

[Unit]
Description=ComRoster — tableau d'affectation beltpacks intercom
# Démarrage au plus tôt : l'affichage local (127.0.0.1) n'attend pas le réseau.
# (Derrière Nginx : le proxy retentera tant que gunicorn n'écoute pas encore.)

[Service]
Type=simple
User=comroster
Group=comroster
WorkingDirectory=/opt/comroster
# Clé de session OBLIGATOIRE : générer avec `python -c "import secrets;print(secrets.token_hex(32))"`
Environment=FLASK_SECRET_KEY=__REMPLACER_PAR_UNE_CLE_LONGUE_ET_ALEATOIRE__
Environment=DATA_DIR=/opt/comroster/instance
# Ce déploiement est derrière Nginx (deploy/nginx.conf) : fait confiance à
# X-Forwarded-For pour que le rate-limit du login voie l'IP réelle des clients.
Environment=COMROSTER_BEHIND_PROXY=true
# Optionnel : délai (s) des requêtes vers l'antenne Bolero (défaut 5).
# Environment=COMROSTER_ANTENNA_TIMEOUT=5
ExecStart=/opt/comroster/.venv/bin/gunicorn -c gunicorn.conf.py app:app
Restart=on-failure
RestartSec=3
# Durcissement systemd. PAS de NoNewPrivileges (voir en-tête : casserait sudo).
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/comroster/instance

[Install]
WantedBy=multi-user.target
